# homedirs/i3lock/.config/sv/i3lock/run
#!/bin/sh
exec xss-lock --transfer-sleep-lock -- i3lock --nofork --image ~/.config/X11/lock/lock.png
