# rootdirs/runit-rshared-root/etc/sv/make-rshared/run
#!/bin/sh
exec 2>&1
set -e
mount --make-rshared /
exec chpst -b make-rshared pause

