# homedirs/usual/.config/sv/xfce-polkit/run
#!/bin/sh
exec /usr/libexec/xfce-polkit 
