# rootdirs/no-usb-wakeup/etc/udev/rules.d/99-no-usb-wakeup.rules
ACTION=="add", SUBSYSTEM=="usb", TEST=="power/control", RUN+="/bin/sh -c 'echo disabled > /sys$devpath/power/wakeup'"

