# homedirs/docker-rootless/.config/sv/docker/run
#!/bin/sh

if command -v dockerd-rootless
then
	exec dockerd-rootless
elif command -v dockerd-rootless.sh
then
	exec dockerd-rootless.sh
fi



